# hospital-pipeline/crontab
* * * * * cd /hospital_cronjob && /usr/local/bin/python hospital_pipeline.py >> cron.log 2>&1
